# docs/attachments/nedi-monitor.service
[Unit]
Description=NeDi monitor daemon
Requires=mariadb.service NetworkManager.service 
# Before=crond.service
After=mariadb.service NetworkManager.service 

[Service]
Type=simple
RemainAfterExit=yes
ExecStart=/var/nedi/moni.pl

[Install]
WantedBy=basic.target
